# deploy/it2.container
# systemd Quadlet unit for the MarianMT en→hi service on the Pi.
#
# Install (rootless, as the Pi user):
#   1. Replace <registry> below with your registry/user.
#   2. mkdir -p ~/.config/containers/systemd
#   3. cp deploy/it2.container ~/.config/containers/systemd/it2.container
#   4. loginctl enable-linger $USER
#   5. systemctl --user daemon-reload && systemctl --user start it2

[Unit]
Description=MarianMT en-hi translator (CTranslate2 int8)

[Container]
Image=<registry>/it2-ct2:latest
PublishPort=8080:8080
Environment=OMP_NUM_THREADS=2
Environment=CT2_INTER_THREADS=1
Environment=CT2_INTRA_THREADS=2
PodmanArgs=--cpus 2 --memory 1500m

[Service]
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=default.target
